# extras/magic-mirror.service
# magic-mirror
#
# Provides a server and frontend for a magic mirror

[Service]
ExecStart=/usr/bin/node server.js
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=magic-mirror
User=pi
Group=pi
Environment=NODE_ENV=production
Environment=MAGIC_MIRROR_CONFIG=/etc/default/magic-mirror
WorkingDirectory=/usr/lib/node_modules/magic-mirror

[Install]
WantedBy=multi-user.target
